# zoot.socket
[Unit]
Description=The Zoot Skoot
PartOf=zoot@.service

[Socket]
ListenStream=0.0.0.0:9999
Accept=Yes

[Install]
WantedBy=sockets.target

